# systemdunits/t212-ghostfolio.timer
[Unit]
Description=Run Trading212 → Ghostfolio import daily

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
